# artifacts/config/autostart/gnome-keyring-secrets.desktop
[Desktop Entry]
Type=Application
Name=Secrets Storage Service
Exec=gnome-keyring-daemon --start --components=secrets
