# artifacts/config/autostart/wayvnc.desktop
[Desktop Entry]
Type=Application
Name=Wayland VNC Server
TryExec=wayvnc
Exec=wayvnc 0.0.0.0
