# artifacts/config/autostart/wl-clip-persist.desktop
[Desktop Entry]
Type=Application
Name=Wayland Clipboard Persistence Service
Exec=wl-clip-persist --clipboard regular
